# deploy/cloudflared.service
[Unit]
Description=Cloudflare Tunnel for vigil (vault.<domain>)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cloudflared
ExecStart=/usr/local/bin/cloudflared tunnel run vigil
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
